# rootfs/etc/s6/registra/run
#!/bin/bash

if [ -f ./setup ]
then
  ./setup
fi

STARTCMD="/bin/registrator"


exec ${STARTCMD} $backend
